# service_examples/my-hello.service
[Unit]
Description = hello daemon

[Service]
ExecStart = /home/dev/service_utils/command_exmpales/hello.sh
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target
